# xulpymoney/mystocks.initd
#!/sbin/runscript

depend() {
	need net
	need postgresql
}

start() {
	ebegin "Starting myquotes"
	start-stop-daemon --start --quiet --exec /usr/bin/mqd&
	eend $?
}


stop() {
	ebegin "Stoping myquotes"
#	killall -9 myquotesd
	pkill -9 -f mq.*.py
#	start-stop-daemon --stop --quiet --exec /usr/bin/myquotesd&
	eend $?
}
